# base/supervisor/supervisor-chromium.ini
[program:chromium]
directory=%HOME%
command=chromium-browser --remote-debugging-port=9222 --no-first-run --no-default-browser-check --user-data-dir=remote-profile --disable-dev-shm-usage --start-maximized
user=%USER%
environment=DISPLAY=":1",HOME="%HOME%",USER="%USER%"